# scripts/goose-relay.service
[Unit]
Description=GooseRelayVPN exit server
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/root
ExecStart=/root/goose-server-linux -config /root/server_config.json
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
